# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# MariaDB Service
[program:mariadb]
command=/usr/sbin/mysqld --user=mysql --pid-file=/run/mysqld/mysqld.pid --socket=/run/mysqld/mysqld.sock --datadir=/var/lib/mysql
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/mariadb.log
redirect_stderr=true

# Wait for MariaDB to be ready before starting other services
[program:wait_mysql]
command=/bin/bash -c 'for i in {1..30}; do mysqladmin ping -h localhost --password=hrs_admin_router -u root > /dev/null 2>&1 && exit 0; sleep 1; done; exit 1'
autostart=true
autorestart=false
priority=2
stdout_logfile=/var/log/supervisor/wait_mysql.log
redirect_stderr=true

# HAProxy Service
[program:haproxy]
command=/bin/bash -c 'envsubst < /etc/haproxy/haproxy.cfg.template > /etc/haproxy/haproxy.cfg && haproxy -f /etc/haproxy/haproxy.cfg'
autostart=true
autorestart=true
priority=3
environment=SERVER_HOSTNAME="127.0.0.1",SERVER_PORT="80"
stdout_logfile=/var/log/supervisor/haproxy.log
redirect_stderr=true

# Mitmproxy Service
[program:mitmproxy]
command=/usr/local/bin/mitmdump --mode reverse:http://127.0.0.1:80 -p 8002 --set block_global=false --no-http2
autostart=true
autorestart=true
priority=4
stdout_logfile=/var/log/supervisor/mitmproxy.log
redirect_stderr=true

# Apache (Flask App) Service
[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/apache2.log
redirect_stderr=true
